# config/thora.bot.service
[Unit]
Description=Thora Bot
After=syslog.target

[Service]
Type=simple
Restart=always
WorkingDirectory=/thora/bot/
EnvironmentFile=/thora/bot/.env
ExecStart=/thora/target/release/bot
KillSignal=SIGQUIT
StandardError=file:/thora/logs/bot.log
NotifyAccess=all

[Install]
WantedBy=multi-user.target
